# scraper/config.cron
# add this line to `crontab -e`
* * * * * cd /home/timon/Code/scraper && /home/timon/.cache/pypoetry/virtualenvs/scraper-omaoyfXU-py3.9/bin/scrapy crawl volume -o /home/timon/Code/scraper/data/volumes.jl > /home/timon/log.txt > /dev/null 2>&1
